# factory/fastapi_app.service
[Unit]
Description=FastAPI Watchdog Service
After=network.target

[Service]
Type=simple
WorkingDirectory={WORKDIR}
ExecStart={VENV}/bin/python -m uvicorn src.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

# PID lock file prevents multiple instances
Environment=SKIP_PID_LOCK=0

[Install]
WantedBy=default.target
